# confirm_sd_service/confirm_sd.timer
[Unit]
Description=Read confirmation of all SD card sectors.

[Timer]
OnBootSec=15min
OnUnitActiveSec=1month

[Install]
WantedBy=timers.target
